# .config/systemd/user/timerexample/dotice.service
[Unit]
Description=Run dotice service

[Service]
Type=oneshot
ExecStart=/home/ztclog/.config/systemd/user/dotice
